# server/supervisor/production.ini
[program:harvest-site-production]
directory = /home/thedevel/sites/harvest-site-env/harvest-site
command = /home/thedevel/sites/harvest-site-env/bin/uwsgi --ini /home/thedevel/sites/harvest-site-env/harvest-site/server/uwsgi/production.ini
autostart = true
autorestart = true
